# docker/config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:layer1]
command=/app/bin/ifr_socket_server
environment=SOCKET_PATH="/tmp/mfn_test_layer1.sock",MFN_DATA_DIR="/data/mfn/memory"
autostart=true
autorestart=true
startretries=5
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:layer2]
command=/app/bin/layer2_socket_server
environment=SOCKET_PATH="/tmp/mfn_test_layer2.sock",MFN_DATA_DIR="/data/mfn/memory",RUST_LOG="info"
autostart=true
autorestart=true
startretries=5
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:layer3]
command=/app/bin/layer3_alm
environment=SOCKET_PATH="/tmp/mfn_test_layer3.sock",MFN_DATA_DIR="/data/mfn/memory"
autostart=true
autorestart=true
startretries=5
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:layer4]
command=/app/bin/layer4_socket_server
environment=SOCKET_PATH="/tmp/mfn_test_layer4.sock",MFN_DATA_DIR="/data/mfn/memory",RUST_LOG="info"
autostart=true
autorestart=true
startretries=5
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:layer5]
command=/app/bin/layer5_socket_server
environment=SOCKET_PATH="/tmp/mfn_test_layer5.sock",MFN_DATA_DIR="/data/mfn/memory",RUST_LOG="info"
autostart=true
autorestart=true
startretries=5
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:gateway]
command=/app/bin/mfn-gateway
environment=MFN_API_PORT="8080",MFN_DATA_DIR="/data/mfn/memory",RUST_LOG="info"
autostart=true
autorestart=true
startretries=5
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
