# docker/config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
pidfile=/app/logs/supervisord.pid
childlogdir=/app/logs
user=root

[unix_http_server]
file=/app/sockets/supervisor.sock
chmod=0700
chown=mfn:mfn

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///app/sockets/supervisor.sock

# MFN Orchestrator - Central coordinator
[program:mfn_orchestrator]
command=/app/scripts/start_orchestrator.sh
user=mfn
autostart=true
autorestart=true
startretries=3
startsecs=10
priority=100
stdout_logfile=/app/logs/orchestrator.log
stderr_logfile=/app/logs/orchestrator_error.log
environment=HOME="/app",USER="mfn",PATH="/app/venv/bin:/usr/local/bin:/usr/bin:/bin"

# Layer 1 - Zig IFR
[program:layer1_ifr]
command=/app/bin/layer1_socket_server
user=mfn
autostart=true
autorestart=true
startretries=3
startsecs=5
priority=200
stdout_logfile=/app/logs/layer1.log
stderr_logfile=/app/logs/layer1_error.log
environment=SOCKET_PATH="/app/sockets/layer1.sock"

# Layer 2 - Rust DSR
[program:layer2_dsr]
command=/app/bin/layer2_socket_server
user=mfn
autostart=true
autorestart=true
startretries=3
startsecs=5
priority=200
stdout_logfile=/app/logs/layer2.log
stderr_logfile=/app/logs/layer2_error.log
environment=SOCKET_PATH="/app/sockets/layer2.sock",RUST_LOG="info"

# Layer 3 - Go ALM
[program:layer3_alm]
command=/app/bin/layer3_server
user=mfn
autostart=true
autorestart=true
startretries=3
startsecs=5
priority=200
stdout_logfile=/app/logs/layer3.log
stderr_logfile=/app/logs/layer3_error.log
environment=SOCKET_PATH="/app/sockets/layer3.sock"

# Layer 4 - Rust CPE
[program:layer4_cpe]
command=/app/bin/layer4_socket_server
user=mfn
autostart=true
autorestart=true
startretries=3
startsecs=5
priority=200
stdout_logfile=/app/logs/layer4.log
stderr_logfile=/app/logs/layer4_error.log
environment=SOCKET_PATH="/app/sockets/layer4.sock",RUST_LOG="info"

# MFN API Gateway
[program:mfn_api]
command=/app/venv/bin/python /app/scripts/api_gateway.py
user=mfn
autostart=true
autorestart=true
startretries=3
startsecs=10
priority=300
stdout_logfile=/app/logs/api.log
stderr_logfile=/app/logs/api_error.log
environment=HOME="/app",USER="mfn",PATH="/app/venv/bin:/usr/local/bin:/usr/bin:/bin",PYTHONPATH="/app/lib"

# MFN Dashboard
[program:mfn_dashboard]
command=/app/venv/bin/python /app/scripts/dashboard_server.py
user=mfn
autostart=true
autorestart=true
startretries=3
startsecs=10
priority=400
stdout_logfile=/app/logs/dashboard.log
stderr_logfile=/app/logs/dashboard_error.log
environment=HOME="/app",USER="mfn",PATH="/app/venv/bin:/usr/local/bin:/usr/bin:/bin",PYTHONPATH="/app/lib"

# Health Monitor
[program:health_monitor]
command=/app/scripts/health_monitor.sh
user=mfn
autostart=true
autorestart=true
startretries=3
startsecs=30
priority=500
stdout_logfile=/app/logs/health_monitor.log
stderr_logfile=/app/logs/health_monitor_error.log

# Persistence Manager
[program:persistence_manager]
command=/app/venv/bin/python /app/scripts/persistence_daemon.py
user=mfn
autostart=true
autorestart=true
startretries=3
startsecs=10
priority=150
stdout_logfile=/app/logs/persistence.log
stderr_logfile=/app/logs/persistence_error.log
environment=HOME="/app",USER="mfn",PATH="/app/venv/bin:/usr/local/bin:/usr/bin:/bin",PYTHONPATH="/app/lib"

# Group definitions for coordinated restart
[group:mfn_layers]
programs=layer1_ifr,layer2_dsr,layer3_alm,layer4_cpe
priority=200

[group:mfn_services]
programs=mfn_orchestrator,mfn_api,mfn_dashboard,health_monitor,persistence_manager
priority=300